# src/model/configs/supervisord.ini
[supervisord]
logfile = /bertynet/supervisord.log
nodaemon = true

[program:flask]
numprocs = 1
numprocs_start = 1
process_name = flask_%(process_num)s
redirect_stderr = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

command=bertynet --port=808%(process_num)s --model_path /bertynet/data/best_model.pt --data_dir /bertynet/data

autostart=true
autorestart=true
